# root/etc/s6-overlay/s6-rc.d/svc-qbit-manager/run
#!/usr/bin/with-contenv bash

# qBittorrent Manager Service
# Runs the HTTP API orchestrator for managing torrent processing and space management

exec s6-setuidgid abc python3 /app/main.py
